# tests/contrib/jujugui/config/haproxy.conf
description "HAProxy"
author "Canonical"

start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]

env CONF=/etc/haproxy/haproxy.cfg
env DAEMON=/usr/sbin/haproxy

expect fork
respawn
respawn limit 10 5

pre-start script
    # Test configuration, exit if errors are found.
    $DAEMON -c -f $CONF
    if [ $? -ne 0 ]
        then exit $?
    fi
end script

exec $DAEMON -f $CONF
